# infrastructure/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=20MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/apex

[program:apex-api]
command=python -m uvicorn trading_system.main:app --host 0.0.0.0 --port 8000 --workers 1 --loop asyncio --log-level info --access-log --no-use-colors
directory=/app
user=apex
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app",HOME="/home/apex"
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stopsignal=SIGTERM
redirect_stderr=false
stdout_logfile=/app/logs/apex-api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/apex-api-err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
priority=10

[program:apex-dashboard]
command=python -m streamlit run trading_system/dashboard/app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true --server.enableCORS=false --server.enableXsrfProtection=false --browser.gatherUsageStats=false --logger.level=info
directory=/app
user=apex
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app",HOME="/home/apex"
autostart=true
autorestart=true
startsecs=15
startretries=3
stopwaitsecs=15
stopsignal=SIGTERM
redirect_stderr=false
stdout_logfile=/app/logs/dashboard.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/dashboard-err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=20

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
stopsignal=QUIT
redirect_stderr=false
stdout_logfile=/app/logs/nginx-access.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/nginx-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
priority=5

[group:apex]
programs=nginx,apex-api,apex-dashboard
priority=999

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
chown=root:root

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface
