# infrastructure/supervisord.conf
; APEX Trading System - Supervisord Process Manager Configuration

[supervisord]
logfile=/var/log/apex/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
nodaemon=true
user=apex

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ---------------------------------------------------------
; FastAPI REST API server (uvicorn)
; ---------------------------------------------------------
[program:api]
command=python -m uvicorn trading_system.api.server:app --host 0.0.0.0 --port 8000 --workers 2 --log-level info
directory=/app
user=apex
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
killasgroup=true
stdout_logfile=/var/log/apex/api.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/apex/api_error.log
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"

; ---------------------------------------------------------
; Streamlit dashboard
; ---------------------------------------------------------
[program:dashboard]
command=python -m streamlit run trading_system/dashboard/app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true --browser.gatherUsageStats false
directory=/app
user=apex
autostart=true
autorestart=true
startretries=5
startsecs=15
stopwaitsecs=15
stdout_logfile=/var/log/apex/dashboard.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/apex/dashboard_error.log
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"

; ---------------------------------------------------------
; Nginx reverse proxy
; ---------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/apex/nginx.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/apex/nginx_error.log
stderr_logfile_maxbytes=10MB

; ---------------------------------------------------------
; Process group for coordinated start/stop
; ---------------------------------------------------------
[group:apex]
programs=api,dashboard,nginx
priority=100
